# velocity-report.service
[Unit]
Description=Radar Monitor Service
After=network.target

[Service]
Type=simple
User=david
Group=david
ExecStart=/home/david/code/velocity.report/radar
WorkingDirectory=/home/david/code/velocity.report
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=go-sensor

[Install]
WantedBy=multi-user.target
