# image/stage-velocity/03-velocity-config/files/velocity-report.service
[Unit]
Description=Velocity.report radar monitor service
After=network.target
# Only start when the radar USB device is connected.
# The udev rule (99-velocity-report.rules) triggers this service via
# SYSTEMD_WANTS when the device appears. BindsTo stops the service
# cleanly when the device is removed.
BindsTo=dev-velocity\x2dradar.device
After=dev-velocity\x2dradar.device

[Service]
User=velocity
Group=velocity
Type=simple
ExecStart=/usr/local/bin/velocity-report --port /dev/velocity-radar --db-path /var/lib/velocity-report/sensor_data.db --config /opt/velocity-report/config/tuning.defaults.json
WorkingDirectory=/var/lib/velocity-report
Environment="PDF_GENERATOR_DIR=/opt/velocity-report/tools/pdf-generator"
Environment="PDF_GENERATOR_PYTHON=/opt/velocity-report/.venv/bin/python"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=velocity-report

[Install]
WantedBy=multi-user.target
