# image/stage-velocity/03-velocity-config/files/velocity-generate-tls.service
[Unit]
Description=Generate TLS certificates for velocity.local
Before=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/velocity-generate-tls.sh /var/lib/velocity-report/tls
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
