# cmd/deploy/velocity-report.service
[Unit]
Description=Velocity.report radar monitor service
After=network.target

[Service]
# Run as a dedicated, non-privileged user created during deployment (see README)
User=velocity
Group=velocity
Type=simple
# The deployed binary (rename/copy from build output to this path)
# Pass explicit DB path so the process finds the SQLite DB in a stable location
# (the binary still defaults to ./sensor_data.db when no flag/env is provided)
ExecStart=/usr/local/bin/velocity-report --db-path /var/lib/velocity-report/sensor_data.db
# A stable working directory for runtime files; owned by the `velocity` user
WorkingDirectory=/var/lib/velocity-report
# PDF generator configuration
Environment="PDF_GENERATOR_DIR=/opt/velocity-report/tools/pdf-generator"
Environment="PDF_GENERATOR_PYTHON=/opt/velocity-report/.venv/bin/python"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=velocity-report

[Install]
WantedBy=multi-user.target
